# activate.service
[Unit]
Description=Calls shell script to run BluePyServer (rfcomm-serverpy3)

[Service]
Type=oneshot

[Install]
WantedBy=multi-user.target